# code/ia-orquestador/deploy/systemd/ia-orquestador.service
[Unit]
Description=IA Orquestador Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orquestador
Group=orquestador
WorkingDirectory=/opt/ia-orquestador
ExecStart=/opt/ia-orquestador/bin/orchestrator \
  -transport=http \
  -http-addr=0.0.0.0:7438 \
  -db=/var/lib/ia-orquestador/orchestrator.db \
  -memory-url=http://<ia-recuerdo-host>:7438 \
  -project=ia-orquestador
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
